# systemd/crypto-bot.service
[Unit]
Description=AI Crypto Trading Telegram Bot
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart={{WORKING_DIRECTORY}}/venv/bin/python3 agent.py bot
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=crypto-bot

[Install]
WantedBy=multi-user.target
